# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:hermes-agent]
command=/opt/hermes-venv/bin/hermes gateway run
user=workspace
environment=
    HOME="/home/workspace",
    HERMES_HOME="/home/workspace/.hermes",
    PATH="/opt/hermes-venv/bin:/usr/local/bin:/usr/bin:/bin",
    GOOGLE_API_KEY="%(ENV_GOOGLE_API_KEY)s",
    API_SERVER_ENABLED="true",
    API_SERVER_HOST="0.0.0.0"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startsecs=15
priority=1

[program:hermes-dashboard]
command=/opt/hermes-venv/bin/hermes dashboard
user=workspace
environment=
    HOME="/home/workspace",
    HERMES_HOME="/home/workspace/.hermes",
    PATH="/opt/hermes-venv/bin:/usr/local/bin:/usr/bin:/bin",
    GOOGLE_API_KEY="%(ENV_GOOGLE_API_KEY)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startsecs=15
priority=2

[program:hermes-workspace]
command=node --max-old-space-size=2048 /app/server-entry.js
user=workspace
directory=/app
environment=
    HOME="/home/workspace",
    NODE_ENV="production",
    PORT="3000",
    HOST="0.0.0.0",
    HERMES_API_URL="http://127.0.0.1:8642",
    HERMES_DASHBOARD_URL="http://127.0.0.1:9119",
    HERMES_ALLOW_INSECURE_REMOTE="1",
    HERMES_PASSWORD="%(ENV_HERMES_PASSWORD)s",
    GOOGLE_API_KEY="%(ENV_GOOGLE_API_KEY)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startsecs=5
priority=3
