# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:hermes-agent]
command=hermes gateway run
user=workspace
environment=
    HOME="/home/workspace",
    HERMES_UID="10010",
    HERMES_DASHBOARD="1",
    HERMES_DASHBOARD_HOST="0.0.0.0",
    HERMES_DASHBOARD_PORT="9119",
    API_SERVER_ENABLED="true",
    API_SERVER_HOST="0.0.0.0",
    GOOGLE_API_KEY="%(ENV_GOOGLE_API_KEY)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startsecs=10

[program:hermes-workspace]
command=node --max-old-space-size=2048 /app/server-entry.js
user=workspace
directory=/app
environment=
    HOME="/home/workspace",
    NODE_ENV="production",
    PORT="3000",
    HOST="0.0.0.0",
    HERMES_API_URL="http://127.0.0.1:8642",
    HERMES_DASHBOARD_URL="http://127.0.0.1:9119",
    HERMES_PASSWORD="%(ENV_HERMES_PASSWORD)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startsecs=5
